# debian/clearwater-snmp-alarm-agent.upstart
description "Clearwater alarm agent"

start on (starting clearwater-monit)
stop on runlevel [!2345]

respawn
# Sleep for 5 seconds on stop, to avoid respawning too quickly
post-stop exec sleep 5

script
  /usr/share/clearwater/bin/clearwater-snmp-alarm-agent-wrapper
end script

# Be certain that the alarm agent itself is no longer running.
# Otherwise, if the wrapper script is killed, the alarm agent
# won't be, and we'll end up with two copies of it.
post-stop script
  pkill -9 -f -x '/usr/share/clearwater/bin/cw_alarm_agent.*' || true
end script
