# debian/clearwater-snmp-alarm-agent.upstart
description "Clearwater alarm agent"

start on (starting clearwater-monit)
stop on runlevel [!2345]

respawn
# Sleep for 5 seconds on stop, to avoid respawning too quickly
post-stop exec sleep 5

script
  /usr/share/clearwater/bin/clearwater-snmp-alarm-agent-wrapper
end script
